# package/base-files/files/etc/uci-defaults/12_network-generate-ula
# shellcheck shell=sh

[ "$(uci -q get network.globals.ula_prefix)" != "auto" ] && exit 0

r1=$(hexdump -v -n 1 -e '1/1 "%02x"' /dev/urandom)
r2=$(hexdump -v -n 2 -e '2/1 "%02x"' /dev/urandom)
r3=$(hexdump -v -n 2 -e '2/1 "%02x"' /dev/urandom)

uci -q batch <<-EOF >/dev/null
	set network.globals.ula_prefix=fd$r1:$r2:$r3::/48
	commit network
EOF

exit 0
