# package/base-files/files/etc/uci-defaults/13_fix-group-user
# shellcheck shell=sh

# shellcheck source=../../lib/functions.sh
. /lib/functions.sh

grep -sl Require-User /usr/lib/opkg/info/*.control | while read -r file; do
	file="${file##*/}"
	add_group_and_user "${file%.control}"
done

exit 0
