# package/base-files/files/etc/uci-defaults/13_fix-group-user
# shellcheck shell=sh

# shellcheck source=../../lib/functions.sh
. /lib/functions.sh

# Skip if we don't have /usr/lib/opkg/info (APK installation)
[ -d /usr/lib/opkg/info ] || exit 0

grep -sl Require-User /usr/lib/opkg/info/*.control | while read -r file; do
	file="${file##*/}"
	add_group_and_user "${file%.control}"
done

exit 0
